# roles/systemd_user_units/files/systemd/deploy@.service
[Unit]
Description=Deploy %i!

[Service]
WorkingDirectory=%h/deploy/%i
ExecStart=bash -c "\
    [ -x deploy ] && exec ./deploy ||\
    [ -x deploy.sh ] && exec ./deploy.sh ||\
    echo deployment not found && exit 1"
Environment=ANSIBLE_LOCALHOST_WARNING=False
Environment=ANSIBLE_INVENTORY_UNPARSED_WARNING=False
EnvironmentFile=-%h/deploy/%i/env

[Install]
WantedBy=deploy.target
