# roles/systemd_user_units/files/systemd/deploy@.path
[Unit]
Description=Deploy %i when deployment files change!

[Path]
PathChanged=%h/deploy/%i/deploy.sh
PathChanged=%h/deploy/%i/deploy.env
TriggerLimitIntervalSec=60
TriggerLimitBurst=3
Unit=deploy@%i.service
