# roles/systemd_user_units/files/systemd/deploy@.service
[Unit]
Description=Deploy %i!

[Service]
WorkingDirectory=%h/deploy/%i
ExecStart=%h/deploy/%i/deploy.sh
Environment=ANSIBLE_LOCALHOST_WARNING=False
Environment=ANSIBLE_INVENTORY_UNPARSED_WARNING=False
EnvironmentFile=-%h/deploy/%i/deploy.env

[Install]
WantedBy=deploy.target
